# nzbmegasearch/services/megasearch/run
#!/bin/bash

exec /sbin/setuser abc python /opt/megasearch/NZBmegasearch/mega2.py
