# nzbmegasearch/services/tail/run
#!/bin/bash

exec /sbin/setuser abc tail -f /opt/megasearch/NZBmegasearch/logs/nzbmegasearch.log
